# trmnl.service
[Unit]
Description=TRMNL display client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/trmnl/trmnl.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
